# DEBIAN/postrm
#!/bin/bash

if [ -f /etc/apache2/conf-enabled/scalix_autodiscover_directory.conf ]; then
    rm -f /etc/apache2/conf-enabled/scalix_autodiscover_directory.conf
fi

invoke-rc.d apache2 force-reload